# audit.service
[Unit]
Description=Meu Script de Inicialização

[Service]
ExecStart=/.FAB/audit.sh
Restart=always

[Install]
WantedBy=graphical.target
