# backup.service
[Unit]
Description=Meu Script de Inicialização

[Service]
ExecStart=/.FAB/backup.sh
Restart=always

[Install]
WantedBy=graphical.target
